# docs/supervisor.conf
; 服务启动的配置
[program:oss]
command=/root/yw/backend/oss_srv
directory=/root/yw/backend
autostart=true
autorestart=true
redirect_stderr=true
startsec=3
priority=10
